# scripts/com.token-usage-dashboard.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.token-usage-dashboard.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>{USER_HOME}/.local/bin/node</string>
    <string>{USER_HOME}/Documents/kimi/token-usage-dashboard/server.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>{USER_HOME}/Documents/kimi/token-usage-dashboard</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/token-usage-dashboard.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/token-usage-dashboard.err.log</string>
</dict>
</plist>
